# install/templates/multivpn-rotate-slot@.timer.tmpl
[Unit]
Description=Rotate multivpn slot %i daily with 12h jitter

[Timer]
OnCalendar=daily
RandomizedDelaySec=12h
Persistent=true
Unit=multivpn-rotate-slot@%i.service

[Install]
WantedBy=timers.target
